# nextcloud/9.0/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=crond -f

[program:php-fpm]
command=php-fpm7

[program:nginx]
command=nginx
