# unmaintained/rutorrent/rootfs/etc/supervisor.d/supervisord.ini
[supervisord]
nodaemon=true

[program:rtorrent]
user=torrent
directory=/home/torrent
priority=2
redirect_stderr=true
environment=HOME="/home/torrent",PWD="/home/torrent",LOGNAME="rtorrent",USER="torrent",TERM="xterm"
startsecs=5
command=rtorrent

[program:php-fpm]
command=php-fpm7 --nodaemonize

[program:nginx]
command=nginx
