# owncloud/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f

[program:php-fpm]
command=php-fpm7.0 --nodaemonize

[program:nginx]
command=nginx
