# zerobin/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm7 --nodaemonize

[program:nginx]
command=nginx
