# rtorrent-flood/supervisord.ini
[supervisord]
nodaemon=true

[program:rtorrent]
user=torrent
directory=/home/torrent
priority=2
redirect_stderr=true
environment=HOME="/home/torrent",PWD="/home/torrent",LOGNAME="rtorrent",USER="torrent",TERM="xterm"
command=rtorrent
autorestart=true

[program:flood]
user=torrent
directory=/usr/flood
redirect_stderr=true
command=/usr/flood/server/bin/www
autorestart=true
